# scripts/systemd/charging-aiot-sync-camera-config.service
[Unit]
Description=Sync charging-aiot camera path config from database to MediaMTX
After=network-online.target mariadb.service mysql.service mediamtx.service
Wants=network-online.target mediamtx.service

[Service]
Type=oneshot
User=www-data
Group=www-data
WorkingDirectory=/var/www/charging-platform/charging-aiot-php
Environment=MEDIAMTX_AUTO_SYNC_RUNTIME_DIR=/tmp/charging-aiot-runtime
ExecStartPre=/usr/bin/mkdir -p /tmp/charging-aiot-runtime
ExecStartPre=/usr/bin/chown www-data:www-data /tmp/charging-aiot-runtime
ExecStart=/usr/bin/php /var/www/charging-platform/charging-aiot-php/scripts/sync_mediamtx_from_db.php
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
